# webserver/debian/hokku-server/lib/systemd/system/hokku-server.service
[Unit]
Description=Hokku Spectra 6 image server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/hokku-server
Environment="PYTHONUNBUFFERED=1"
Environment="HOKKU_CONFIG=/etc/hokku/config.json"
DynamicUser=yes
StateDirectory=hokku
CacheDirectory=hokku
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
